# masto-bsky.service
[Unit]
Description=Run masto-bsky script

[Service]
User=nvm
ExecStart=/home/nvm/masto-bsky/main.sh

[Install]
WantedBy=default.target
