# masto-bsky.service
[Unit]
Description=Run masto-bsky script

[Service]
User=user
ExecStart=/home/user/masto-bsky/main.sh

[Install]
WantedBy=default.target
