# masto-bsky.service
[Unit]
Description=Run masto-bsky script

[Service]
User=mastodon
ExecStart=/home/mastodon/masto-bsky/main.sh

[Install]
WantedBy=default.target
